# zenohd.service
[Unit]
Description=Zenoh Router Daemon
Documentation=https://zenoh.io/docs/
After=network.target

[Service]
Type=simple
# Try full-featured zenohd first, fall back to basic
ExecStart=bash -c "if [ -x /usr/local/bin/zenohd-full ]; then exec /usr/local/bin/zenohd-full --listen tcp/[::]:7447 --rest-http-port 7447; else exec %h/.cargo/bin/zenohd --listen tcp/[::]:7447; fi"
Restart=on-failure
RestartSec=5

# Environment variables (optional)
Environment=ZENOH_TASK_PRIORITY=0
# Uncomment to persist across reboots
# RemainAfterExit=yes

[Install]
WantedBy=default.target
