# zenohd.service
[Unit]
Description=Zenoh Router Daemon
Documentation=https://zenoh.io/docs/
After=network.target

[Service]
Type=simple
# Basic Zenoh router for FlatBuffers native protocol
ExecStart=%h/.cargo/bin/zenohd --listen tcp/[::]:7447
Restart=on-failure
RestartSec=5

# Environment variables (optional)
Environment=ZENOH_TASK_PRIORITY=0
# Uncomment to persist across reboots
# RemainAfterExit=yes

[Install]
WantedBy=default.target
